# autostart/pavutray.desktop
[Desktop Entry]
Type=Application
Exec=pavutray
TryExec=pavutray

Name=Pavucontrol-qt tray
